# server_setup/expmonitor.service
[Unit]
Description=Experiment Monitoring Software
After=multi-user.target


[Service]
Environment=PYTHONPATH=/root/experiment-monitoring/src
ExecStart=/usr/bin/python3 /root/experiment-monitoring/src/expmonitor/exec.py
Restart=always
RestartSec=15s   


[Install]
WantedBy=multi-user.target
